# deploy/systemd/jarvis-api.service
[Unit]
Description=Jarvis API
After=network.target

[Service]
Type=simple
User=justin
WorkingDirectory=/home/justin/jarvis2
EnvironmentFile=/home/justin/jarvis2/deploy/.env.prod
ExecStart=/home/justin/.local/bin/uv run uvicorn jarvis.main:app --host 127.0.0.1 --port 8000 --app-dir src
Restart=always
RestartSec=5
WatchdogSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
