# deploy/systemd/jarvis-scheduler.service
[Unit]
Description=Jarvis Scheduler Tick
After=network.target jarvis-api.service

[Service]
Type=simple
User=justin
WorkingDirectory=/home/justin/jarvis2
EnvironmentFile=/home/justin/jarvis2/deploy/.env.prod
ExecStart=/bin/bash -lc 'while true; do /home/justin/.local/bin/uv run python -c "from jarvis.tasks.scheduler import scheduler_tick; scheduler_tick()"; sleep 30; done'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
