# deploy/systemd/jarvis-worker.service
[Unit]
Description=Jarvis Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/agent-framework
EnvironmentFile=/etc/jarvis/jarvis.env
ExecStart=/usr/bin/env uv run celery -A jarvis.celery_app worker -Q agent_priority,agent_default,tools_io,local_llm --loglevel=info
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
